# supervisord.conf
[supervisord]
nodaemon=true
root=true

[program:X11]
command=Xvfb :0 -screen 0 "%(ENV_DISPLAY_WIDTH)s"x"%(ENV_DISPLAY_HEIGHT)s"x24
autorestart=true

[program:x11vnc]
command=/usr/bin/x11vnc -rfbauth /root/.vnc/passwd
autorestart=true

[program:novnc]
command=/root/noVNC/utils/novnc_proxy --vnc localhost:5900 --listen 8083
autorestart=true

[program:socat]
command=socat tcp-listen:6000,reuseaddr,fork unix:/tmp/.X11-unix/X0
autorestart=true

[program:fluxbox]
command=fluxbox
autorestart=true

[program:tsugite]
command=sh -c "cd /code && python3 Tsugite_app.py"
autorestart=true
